# crontab
0 3 * * * root docker exec l2vibe_mongo_service /backup.sh